# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = cron -f -L 15 > /var/log/dockeroutput.log
autostart=true
autorestart=true
